# tiny-setup-package/install-package/exporter/mysql/mysql_exporter.service
[Unit]
Description=Prometheus

[Service]
Environment=
ExecStart=
Restart=always

[Install]
WantedBy=multi-user.target